# base/server/share/lib/systemd/system/pki-quarkusd@.service
# Systemd service template for Quarkus-based PKI subsystems
#
# Each PKI subsystem runs as an independent Quarkus process.
# Instance name format: <instance>-<subsystem>
# Example: pki-quarkus-ca, pki-quarkus-kra, pki-quarkus-ocsp
#
# Usage:
#   systemctl start pki-quarkusd@pki-quarkus-ca
#   systemctl stop pki-quarkusd@pki-quarkus-ca
#   systemctl status pki-quarkusd@pki-quarkus-ca

[Unit]
Description=Dogtag PKI Quarkus Subsystem (%i)
After=network.target
After=syslog.target
After=dirsrv.target

[Service]
Type=simple
User=pkiuser
Group=pkiuser

# Working directory
WorkingDirectory=/var/lib/pki

# Execute via startup script which parses instance/subsystem from %i
ExecStart=/usr/share/pki/scripts/pki-quarkus-start %i

# Graceful shutdown
TimeoutStopSec=30

# Restart on failure
Restart=on-failure
RestartSec=10

# Security hardening
# NOTE: Cannot use ProtectSystem=full because /var/lib/pki/<instance>/conf
# is a symlink to /etc/pki/<instance>, and NSS requires write access to
# the database files under conf/alias.
ProtectHome=true
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/etc/pki /var/lib/pki /var/log/pki

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pki-quarkusd-%i

[Install]
WantedBy=pki-quarkusd.target
